# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:laravel-setup]
command=/bin/bash -c "sleep 15 && cp .env.dokploy .env && echo 'Environment file copied' && php artisan config:clear && php artisan cache:clear && php artisan migrate --force && echo 'Migrations completed' && php artisan db:seed --force && echo 'Seeding completed' && php artisan config:cache && php artisan route:cache && php artisan view:cache && echo 'Laravel caching completed'"
autostart=true
autorestart=false
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
